# recipes-kovan/kovan-kmod/kovan-kmod/kovan.service
[Unit]
Description=kovan-kmod

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=insmod /lib/modules/`uname -r`/drivers/kovan/kovan.ko
